# systemd/metasiberia-map-progress.timer
[Unit]
Description=Run Metasiberia map progress sampler every 15 minutes

[Timer]
OnCalendar=*:0/15
RandomizedDelaySec=120
Persistent=true
Unit=metasiberia-map-progress.service

[Install]
WantedBy=timers.target

