# systemd/metasiberia-map-progress.service
[Unit]
Description=Metasiberia map progress sampler
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/metasiberia_map_maintenance.py sample
User=root

